# root/etc/uci-defaults/99-autoupdate
#!/bin/sh

[ -f /etc/config/autoupdate ] || {
    touch /etc/config/autoupdate
    uci add autoupdate settings
    uci commit autoupdate
}

exit 0

